# rbackups/rbackups.cron
# Every Saturday at 1:15
15 1 *  *  sat  root  /srv/rbackups/rbackups.sh
